# SVAGA/monitoring/load-alert.timer
[Unit]
Description=Каждые 2 мин — мониторинг нагрузки (Telegram)

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
