# SVAGA/monitoring/post-reboot-check.service
[Unit]
Description=После reboot: проверить VPN/monitoring и уведомить в Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hysteria-monitor
ExecStart=/opt/hysteria-monitor/bin/post-reboot-check.sh

[Install]
WantedBy=multi-user.target

